# crontab
55 6-22 * * * cd ~/FutuOpenDClient/FeedFilter; ./index.js
